# packaging/linux/udev/99-qmk-window-notifier.rules.template
# Allow access to hidraw devices
KERNEL=="hidraw*", ATTRS{idVendor}=="{vendor_id}", ATTRS{idProduct}=="{product_id}", MODE="0666"
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="{vendor_id}", ATTRS{idProduct}=="{product_id}", \
  RUN+="/usr/bin/su {username} -c 'XDG_RUNTIME_DIR=/run/user/{user_id} systemctl --user start qmk-window-notifier.service'"
ACTION=="remove", SUBSYSTEMS=="usb", ATTRS{idVendor}=="{vendor_id}", ATTRS{idProduct}=="{product_id}", \
  RUN+="/usr/bin/su {username} -c 'XDG_RUNTIME_DIR=/run/user/{user_id} systemctl --user stop qmk-window-notifier.service'"
